# config/aios.service
[Unit]
Description=Billiam OS — FOSS AI Operating System Orchestration Layer Daemon
Documentation=https://github.com/iknowkungfubar/billiam-os
After=network.target sound.target

[Service]
Type=simple
Environment=BILLIAM_API_BASE=http://localhost:8080/v1
Environment=BILLIAM_MODEL=qwen-2.5-coder-3b-instruct
Environment=BILLIAM_MEMORY_PATH=%h/.config/aios/memory.json
Environment=BILLIAM_TEMPERATURE=0.2
Environment=BILLIAM_MAX_TOKENS=512
ExecStart=%h/billiam-os/.venv/bin/python -m core.ai_core
WorkingDirectory=%h/billiam-os
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
